# lavandula/systemd/990-nightly.timer
[Unit]
Description=Nightly 990 index refresh and auto-process

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
